# cron/crontab
1 0 * * * pg_dump -h housemate_db -U postgres postgres > /home/postgres/bak/$(date +%Y%m%d%H%M)postgres.bak
